# Chapter07/sample5/calculator.sh
#!/bin/bash

### BEGIN INIT INFO
# Provides:          calculator
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Calculator application
### END INIT INFO

java -jar /var/calculator/calculator.jar &